# board/common/overlay/data/etc/crontabs/root
*/30 * * * * /bin/photoframe.sh sync
0    7 * * * /bin/photoframe.sh display on
0   20 * * * /bin/photoframe.sh display off
